# deployments/installer/linux/surgewave-control.service
[Unit]
Description=Storm Control UI - Web-based administration dashboard
Documentation=https://github.com/Kuestenlogik/Storm
After=storm-broker.service
Wants=storm-broker.service

[Service]
Type=exec
ExecStart=/opt/storm/control/Kuestenlogik.Storm.Control --urls http://0.0.0.0:5050
WorkingDirectory=/opt/storm

# Environment
Environment=DOTNET_RUNNING_IN_CONTAINER=false
Environment=ASPNETCORE_ENVIRONMENT=Production

# Restart policy
Restart=always
RestartSec=10
StartLimitIntervalSec=60
StartLimitBurst=5

# Security
User=storm
Group=storm
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=storm-control

[Install]
WantedBy=multi-user.target
